# packaging/ubuntu-trusty/build-inside/upstart/ourcms.conf
description "ourcms"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
chdir /opt/ourcms/cms
exec ../bin/uwsgi --master --processes 4 --die-on-term --socket :3031 --wsgi-file ourcms/wsgi.py
